# backend/systemd/shareverse-refunds.timer
[Unit]
Description=Run ShareVerse expired buy-together refund processor every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=shareverse-refunds.service

[Install]
WantedBy=timers.target
